# deploy/quadlet/borgmatic-manager.container
# Podman Quadlet container unit for borgmatic-manager
# Install to: ~/.config/containers/systemd/ (rootless) or /etc/containers/systemd/ (root)
# Then: systemctl --user daemon-reload && systemctl --user start borgmatic-manager

[Unit]
Description=Borgmatic Manager - automated backup orchestration
After=local-fs.target

[Container]
Image=ghcr.io/lugoues/borgmatic-manager:latest
Volume=/run/podman/podman.sock:/var/run/docker.sock
Volume=%h/.config/borgmatic-manager:/etc/borgmatic-manager:ro
Volume=borgmatic-cache.volume:/root/.cache/borg
Environment=CONFIG_DIR=/etc/borgmatic-manager
# Optional: SSH keys for remote borg repositories
# Volume=%h/.ssh:/root/.ssh:ro
# Optional: override container socket path
# Environment=CONTAINER_SOCKET=/var/run/docker.sock

[Service]
Restart=always

[Install]
WantedBy=default.target
